# scripts/run_scripts.timer
[Unit]
Description=Run Script

[Timer]
OnCalendar=*-*-* *:00/15:05
Persistent=true

[Install]
WantedBy=timers.target